# roles/gcwebarchives/files/gcwebarchives.service
[Unit]
Description=GCWebArchives
After=syslog.target network.target

[Service]
User=lacwayback
ExecStart=/opt/gcwebarchives/gcwebarchives.jar
SuccessExitStatus=143

[Install]
WantedBy=multi-user.target
